# rsi-monitor.service
[Unit]
Description=SOL RSI Monitor v1.3
Documentation=https://github.com/your-username/rsi-monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rsi-monitor
ExecStart=/home/ubuntu/.nvm/versions/node/v18.20.4/bin/node --max-old-space-size=512 src/server.js
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5

# 环境变量
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/rsi-monitor/.env

# 内存限制（超过 512M 自动重启）
MemoryMax=512M
MemorySwapMax=0

# 日志由 systemd journal 接管，stdout/stderr 均收录
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rsi-monitor

# 安全加固（可选，不影响功能）
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
